# config/timers/tars-health-audit.service
[Unit]
Description=T.A.R.S Health Audit
After=tars.service

[Service]
Environment=TARS_HOME=/opt/tars
Type=oneshot
User=tars
Group=tars
WorkingDirectory=/opt/tars
ExecStart=/opt/tars/scripts/health-audit.sh --report
StandardOutput=append:/var/log/tars-security.log
StandardError=append:/var/log/tars-security.log
