# config/timers/tars-memory-context.service
[Unit]
Description=T.A.R.S Memory Context Regen
After=tars.service

[Service]
Environment=TARS_HOME=/opt/tars
Type=oneshot
User=tars
Group=tars
WorkingDirectory=/opt/tars
ExecStart=/opt/tars/scripts/regen-memory-context.sh
StandardOutput=append:/var/log/tars-context.log
StandardError=append:/var/log/tars-context.log
