# config/tars-rescue.service
# Rescue Bot — standalone systemd service (unsandboxed)
# Install: sudo cp config/tars-rescue.service /etc/systemd/system/ && sudo systemctl daemon-reload

[Unit]
Description=T.A.R.S Rescue Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tars
Group=tars
WorkingDirectory=/opt/tars
ExecStart=/usr/local/bin/uv run --no-sync python -m src.main --profile rescue
Restart=always
RestartSec=10
# HOME and PATH are set dynamically by setup.sh based on the service account
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# NO sandbox — Rescue Bot needs full system access for dev/ops tasks

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tars-rescue

[Install]
WantedBy=multi-user.target
