# config/tars-memory-decay.service
[Unit]
Description=T.A.R.S Memory Context Regen + Decay
After=tars.service

[Service]
Type=oneshot
User=tars
Group=tars
WorkingDirectory=/opt/tars
ExecStart=/opt/tars/scripts/memory-decay.sh
