# config/tars.service
# T.A.R.S — systemd service template

[Unit]
Description=T.A.R.S Agent System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tars
Group=tars
WorkingDirectory=/opt/tars
ExecStart=/usr/local/bin/uv run --no-sync python -m src.main
Restart=always
RestartSec=10
# HOME and PATH are set dynamically by setup.sh based on the service account
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateDevices=true
PrivateTmp=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
UMask=0077
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=~@mount @reboot @swap @debug @obsolete @cpu-emulation @raw-io
CapabilityBoundingSet=
LockPersonality=true
RemoveIPC=true
ReadWritePaths=/opt/tars/data /opt/tars/agents /opt/tars/config /tmp %h/.cache %h/.claude
ReadOnlyPaths=/opt/tars

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tars

[Install]
WantedBy=multi-user.target
